# build/daemon/Packerspike/buildDeb/initd
#!/bin/sh

### BEGIN INIT INFO
# Provides:          packer-spike
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Required-Start:
# Required-Stop:
# Description: Control Script for packer-spike
### END INIT INFO


service=packer-spike

daemontools_service=/service/packer-spike


if [ -e /etc/default/packer-spike ]; then
  . /etc/default/packer-spike
fi


stop() {
    touch $daemontools_service/down
    svc -d $daemontools_service
}

start() {
    rm -f $daemontools_service/down
    svc -u $daemontools_service
}

restart() {
    svc -t $daemontools_service
}

case $1 in

start)
    echo -n "${service}: "
    start
    if [ $? -eq 0 ]
    then

        echo "OK"
    else
        echo "ERROR"

    fi
    ;;
 stop)
    echo -n "$service:  "
    stop
    if [ $? -eq 0 ]
    then

        echo "OK"
    else
        echo "ERROR"

    fi
    ;;
restart)
    echo -n "$service: "
    restart
    if [ $? -eq 0 ]
    then

        echo "OK"
    else
        echo "ERROR"

    fi
    ;;
  *)
    echo "Usage: service $service {start|stop|restart}"
    exit 1
esac