# build/debian/postinst
#!/bin/sh -e

ec() {
    echo "$@" >&2
    "$@"
}

case "$1" in
    configure)
        

        
        [ -x /bin/touch ] && touch=/bin/touch || touch=/usr/bin/touch

        
        $touch /service/packer-spike/down

        
        /usr/sbin/update-rc.d packer-spike start 85 2 3 4 5 . stop 15 0 1 6 .

        
        ;;
esac
