# pkg_scripts/spintest.conf
description "spintest"

expect fork

start on runlevel [2345]
stop on runlevel [!2345]

exec java -jar /opt/spintest/spintest*.jar 2>&1 > /var/log/spintest/log.log &
